# scripts/com.flux.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<!-- Label must match the filename (minus .plist) -->
	<key>Label</key>
	<string>com.flux.daemon</string>

	<!-- Working directory for the process -->
	<key>WorkingDirectory</key>
	<string>/path/to/flux</string>

	<!-- Command to run -->
	<key>ProgramArguments</key>
	<array>
		<string>/path/to/bun</string>
		<string>run</string>
		<string>start</string>
	</array>

	<!-- Environment variables (adjust for your deployment) -->
	<key>EnvironmentVariables</key>
	<dict>
		<key>NODE_ENV</key>
		<string>production</string>
		<key>CONVEX_URL</key>
		<string>https://your-deployment.convex.cloud</string>
	</dict>

	<!-- Keep alive: restart if the process exits -->
	<key>KeepAlive</key>
	<true/>

	<!-- Logging -->
	<key>StandardOutPath</key>
	<string>/tmp/flux-daemon.stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/flux-daemon.stderr.log</string>

	<!--
	  ExitTimeOut: seconds launchd waits after SIGTERM before sending SIGKILL.
	  Set to 90s: 60s for active sessions to complete + 30s buffer.
	  The graceful shutdown handler uses a 60s timeout for session draining,
	  so 90s ensures launchd doesn't SIGKILL while we're still cleaning up.
	-->
	<key>ExitTimeOut</key>
	<integer>90</integer>
</dict>
</plist>
